# deploy/ikamdocs.service
[Unit]
Description=ikamdocs FastAPI backend
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/ikamdocs/backend
Environment="PATH=/var/ikamdocs/backend/venv/bin"
ExecStartPre=/bin/bash -c 'pid=$(/usr/bin/lsof -ti :8000 2>/dev/null); [ -n "$pid" ] && /bin/kill -9 $pid 2>/dev/null; /bin/sleep 1'
ExecStart=/var/ikamdocs/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
ExecReload=/bin/kill -s HUP $MAINPID
TimeoutStopSec=10
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
